# systemd/sortinghat.timer
[Unit]
Description=Runs SortingHat on a schedule

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target

